# deploy/supervisor-worker.conf
; Supervisor config for Laravel queue worker
; Path: /etc/supervisor/conf.d/globalafricaplus-worker.conf

[program:globalafricaplus-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/globalafricaplus.com/current/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=90
autostart=true
autorestart=true
user=deploy
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/globalafricaplus.com/shared/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
